# misc/rules.d/70-snap.scrcpy.rules
# This file is automatically generated.
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="03f0", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="03fc", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0408", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0409", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0414", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0451", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0471", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0482", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0489", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="04b7", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="04c5", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="04da", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="04dd", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="04e8", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0502", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0531", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="054c", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="05c6", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="067e", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="091e", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0930", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0955", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0b05", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0bb4", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0c2e", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0db0", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0e79", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0e8d", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0f1c", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0fce", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1004", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="109b", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="10a9", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1219", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="12d1", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1662", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="16d5", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="17ef", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="18d1", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1949", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="19a5", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="19d2", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1b8e", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1bbb", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1d09", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1d45", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1d4d", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1d91", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1e85", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1ebf", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1f3a", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1f53", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2006", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="201e", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2080", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2116", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2207", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2237", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2257", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="22b8", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="22d9", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2314", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2340", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2420", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="24e3", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="25e3", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2717", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="271d", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2836", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2916", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="297f", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="29a9", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="29e4", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2a45", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2a47", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2a49", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2a70", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2ae5", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="413c", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="8087", TAG+="snap_scrcpy_adb"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="e040", TAG+="snap_scrcpy_adb"
# opengl
KERNEL=="nvhost-*", TAG+="snap_scrcpy_adb"
# opengl
KERNEL=="nvmap", TAG+="snap_scrcpy_adb"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_scrcpy_adb"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_scrcpy_adb"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_scrcpy_adb"
# opengl
KERNEL=="vchiq", TAG+="snap_scrcpy_adb"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_scrcpy_adb"
TAG=="snap_scrcpy_adb", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_scrcpy_adb $devpath $major:$minor"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="03f0", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="03fc", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0408", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0409", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0414", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0451", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0471", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0482", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0489", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="04b7", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="04c5", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="04da", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="04dd", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="04e8", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0502", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0531", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="054c", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="05c6", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="067e", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="091e", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0930", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0955", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0b05", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0bb4", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0c2e", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0db0", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0e79", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0e8d", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0f1c", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="0fce", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1004", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="109b", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="10a9", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1219", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="12d1", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1662", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="16d5", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="17ef", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="18d1", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1949", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="19a5", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="19d2", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1b8e", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1bbb", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1d09", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1d45", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1d4d", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1d91", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1e85", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1ebf", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1f3a", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="1f53", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2006", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="201e", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2080", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2116", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2207", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2237", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2257", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="22b8", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="22d9", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2314", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2340", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2420", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="24e3", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="25e3", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2717", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="271d", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2836", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2916", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="297f", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="29a9", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="29e4", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2a45", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2a47", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2a49", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2a70", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="2ae5", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="413c", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="8087", TAG+="snap_scrcpy_scrcpy"
# adb-support
SUBSYSTEM=="usb", ATTR{idVendor}=="e040", TAG+="snap_scrcpy_scrcpy"
# opengl
KERNEL=="nvhost-*", TAG+="snap_scrcpy_scrcpy"
# opengl
KERNEL=="nvmap", TAG+="snap_scrcpy_scrcpy"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_scrcpy_scrcpy"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_scrcpy_scrcpy"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_scrcpy_scrcpy"
# opengl
KERNEL=="vchiq", TAG+="snap_scrcpy_scrcpy"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_scrcpy_scrcpy"
TAG=="snap_scrcpy_scrcpy", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_scrcpy_scrcpy $devpath $major:$minor"
